# etc/metadata.cron
*/5 * * * * root wget -qO- 'http://md.swamid.se/md/swamid-idp-transitive.xml' > /var/run/swamid-idp-transitive.xml
